# config/systemd/user/convert-wallpapers.service
[Unit]
Description=Convert wallpapers for i3-lock

[Service]
Type=oneshot
ExecStart=%h/.scripts/convert_wallpapers.sh
