# config/systemd/user/convert-wallpapers.timer
[Unit]
Description=Convert wallpapers daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target