# config/systemd/user/rss-check.service
[Unit]
Description=Updates RSS count via newsboat

[Service]
Type=oneshot
ExecStart=%h/.scripts/rss/check-updates.sh
